# systemd/os-dim.service
[Unit]
Description=Dim display
After=graphical.target

[Service]
Type=oneshot
ExecStart=/usr/bin/displayctl dim
