# systemd/os-dim.timer
[Timer]
OnIdleSec=8min

[Install]
WantedBy=timers.target
