# systemd/os-screenoff.timer
[Timer]
OnIdleSec=10min

[Install]
WantedBy=timers.target
