# systemd/os-wake.service
[Unit]
Description=Wake display
After=graphical.target

[Service]
Type=oneshot
ExecStart=/usr/bin/displayctl on
